# launchd.com.westlee.openclaw_codex.telegram.bridge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.westlee.openclaw_codex.telegram.bridge</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>/Users/westlee/Projects/openclaw-codex-telegram-bridge/bridge.py</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/westlee/Projects/openclaw-codex-telegram-bridge</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/Users/westlee/.nvm/versions/node/v24.13.1/bin:/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/westlee/Projects/openclaw-codex-telegram-bridge/logs/launchd-bridge.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/westlee/Projects/openclaw-codex-telegram-bridge/logs/launchd-bridge.err.log</string>
</dict>
</plist>
